# debian/cocaine-flow.upstart
description "Cocaine-flow daemon"

start on filesystem
stop on runlevel [!2345]

respawn

script
    # modify these in /etc/default/$UPSTART_JOB (/etc/default/cocaine-flow)
    FLOW=/usr/bin/$UPSTART_JOB
    FLOW_OPTS=
    if [ -f /etc/default/$UPSTART_JOB ]; then
        . /etc/default/$UPSTART_JOB
    fi
    "$FLOW" $FLOW_OPTS
end script
