# upstart/cocaine-flow.conf
# Cocaine-flow backend
# Author: Anton Tyurin  <noxiouz@yandex.ru>
start on runlevel 2
stop on runlevel [!2]
respawn

exec /usr/bin/cocaine-flow --config=/etc/cocaine-flow/config-default.cfg
